# board/laird/som60/rootfs-additions/etc/udev/rules.d/50-emac.rules
ACTION=="add", SUBSYSTEM=="net", DRIVERS=="macb", RUN+="/usr/bin/setmac.sh"
