# board/laird/wb50n_sysd/rootfs-additions/etc/udev/rules.d/50-emac.rules
ACTION=="add", SUBSYSTEM=="net", RUN+="/usr/bin/setmac.sh"
