# package/lrd/lrd-encrypted-storage-toolkit/rootfs/usr/lib/systemd/system/systemd-timedated.service.d/10-timezone-access.conf
[Service]
ReadWritePaths=/data/misc
ReadOnlyPaths=/etc
