# board/laird/wb50n_sysd/rootfs-additions/etc/systemd/system/swupdate.service.d/tmpdir.conf
[Service]
Environment="TMPDIR=/var/tmp"
