# board/laird/som60/rootfs-additions/etc/systemd/system/fskey.service
[Unit]
Description=Filesystem Encryption Key
DefaultDependencies=false

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=keyctl search %:_builtin_fs_keys logon fscrypt:ffffffffffffffff @us
KeyringMode=inherit
